# systemd/templates/systemd/cbsd.target
[Unit]
Description=Clyso Build System service daemon instances

[Install]
WantedBy=multi-user.target
